# files/iusebash.init
#!/bin/sh /etc/rc.common

START=99
STOP=99
USE_PROCD=1

boot() {
	rm -f /tmp/.busybox_ash_history 2>/dev/null
	ln -s /tmp/.busybox_ash_history $HOME/.bash_history
}
